# Roles/roles/mongodb_exporter/templates/mongodb_exporter.service.j2
[Unit]
Description=MongoDB Exporter
User="{{ mongodb_exporter_system_user }}"

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/mongodb_exporter

[Install]
WantedBy=multi-user.target