# Roles/roles/jaeger/files/jaeger-collector.service
[Unit]
Description=Jaeger Collector
After=network.target
Requires=network.target

[Service]
Type=simple
User=jaeger
Environment=SPAN_STORAGE_TYPE=elasticsearch
ExecStart=/usr/local/bin/jaeger-collector --config-file /etc/jaeger/collector.yaml
Restart=always
RestartSec=5
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target