# Roles/roles/jaeger/files/jaeger-query.service
[Unit]
Description=Jaeger Query
After=network.target
Requires=network.target

[Service]
Type=simple
User=jaeger
Environment=SPAN_STORAGE_TYPE=elasticsearch
ExecStart=/usr/local/bin/jaeger-query --config-file /etc/jaeger/query.yaml
Restart=always
RestartSec=5
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target